# systemd/lxc-autoscaler.service
[Unit]
Description=LXC Autoscaler for Proxmox VE
Documentation=https://github.com/example/lxc-autoscaler
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=exec
User=root
Group=root
ExecStart=/usr/local/bin/lxc-autoscaler --config /etc/lxc-autoscaler/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/run /etc/lxc-autoscaler
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Resource limits
LimitNOFILE=65536
TasksMax=512

# Working directory and environment
WorkingDirectory=/var/lib/lxc-autoscaler
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/usr/local/lib/python3/dist-packages

# PID file
PIDFile=/var/run/lxc-autoscaler.pid

# Standard streams
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lxc-autoscaler

[Install]
WantedBy=multi-user.target